# source/systemd/user/sync-repositories.timer
[Unit]
Description=Synchronize repositories

[Timer]
OnBootSec=1min
OnUnitActiveSec=3min
RandomizedDelaySec=60sec

[Install]
WantedBy=timers.target
