# source/systemd/user/sync-news.timer
[Unit]
Description=Synchronize news

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
RandomizedDelaySec=60sec

[Install]
WantedBy=timers.target
