# source/systemd/user/thunderbird-unread.service
[Unit]
Description=Thunderbird Unread

[Service]
Type=fork
ExecStart=%h/.dotfiles/bin/thunderbird-unread-wrapper

[Install]
WantedBy=default.target
