# source/systemd/user/thunderbird-unread.service
[Unit]
Description=Thunderbird Unread

[Service]
Type=simple
ExecStart=%h/bin/dot/thunderbird-unread-wrapper

[Install]
WantedBy=default.target
